# zigbee2mqtt/rootfs/etc/services.d/zigbee2mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Zigbee2mqtt
# Runs Zigbee2mqtt
# ==============================================================================

cd /opt/zigbee2mqtt

exec npm start