# zigbee2mqtt/rootfs/etc/services.d/zigbee2mqtt/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Zigbee2MQTT
# Starts Zigbee2MQTT
# ==============================================================================

# Use share directory for configuration data
export ZIGBEE2MQTT_DATA="/config/zigbee2mqtt"

# Set Node to Production
export NODE_ENV=production

# Optionally Set Zigbee Herdsman Debug
if bashio::config.true 'zigbee_herdsman_debug'; then
    bashio::log.info "Zigbee Herdsman debug logging enabled"
    export DEBUG="zigbee-herdsman:*"
fi

cd /opt/zigbee2mqtt || bashio::exit.nok "Could not change directory to Zigbee2mqtt"

bashio::log.info 'Starting Zigbee2MQTT...'

exec node index.js