# zigbee2mqtt/rootfs/etc/services.d/zigbee2mqtt/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Community Hass.io Add-ons: Zigbee2mqtt
# Runs Zigbee2mqtt
# ==============================================================================
# shellcheck disable=SC1091
source /usr/lib/hassio-addons/base.sh

cd /opt/zigbee2mqtt

exec npm start